# conf/unit/nebula-metad.service
[Service]
Type=simple
ExecStart=/usr/local/nebula/scripts/nebula.service start metad
RemainAfterExit=yes
 
[Install]
WantedBy=multi-user.target
