# conf/unit/nebula.service
[Service]
Type=simple
ExecStart=/usr/local/nebula/scripts/nebula.service start all
RemainAfterExit=yes
 
[Install]
WantedBy=multi-user.target